# mousekeys.service
[Unit]
Description=Mouse Keys daemon for Hyprland
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStartPre=/bin/sh -c 'missing=""; for v in XDG_RUNTIME_DIR WAYLAND_DISPLAY HYPRLAND_INSTANCE_SIGNATURE; do if [ -z "$$(/usr/bin/printenv "$$v")" ]; then missing="$$missing $$v"; fi; done; if [ -n "$$missing" ]; then echo "mousekeys.service requires a Hyprland Wayland session environment; missing:$$missing. Start Hyprland with UWSM or import XDG_RUNTIME_DIR, WAYLAND_DISPLAY, and HYPRLAND_INSTANCE_SIGNATURE into the systemd user manager." >&2; exit 1; fi'
ExecStart=/usr/bin/env mousekeys daemon
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=graphical-session.target
